# deploy/stormtrooper/systemd/rfd-discord-bot-watchdog.service
[Unit]
Description=RFD Discord Bot health watchdog
After=default.target

[Service]
Type=oneshot
WorkingDirectory=/home/paul/agent-work/repos/rfd-discord-bot/deploy/stormtrooper
ExecStart=/home/paul/agent-work/repos/rfd-discord-bot/deploy/stormtrooper/rfd-bot-watchdog.sh
